# systemd/forseti-meshd.service
[Unit]
Description=Forseti Mesh Daemon (forseti-meshd)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/ubuntu/forseti.life/forseti-meshd
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5
Environment="MESHD_BASE_URL=https://forseti.life"
Environment="MESHD_DATABASE_URL=postgresql://forseti_mesh:forseti_mesh_dev@localhost/forseti_mesh"
Environment="MESHD_KEY_DIR=/etc/forseti-meshd/keys"
EnvironmentFile=-/etc/forseti-meshd/config.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/forseti-meshd /tmp

[Install]
WantedBy=multi-user.target
